# chainmart.service
[Unit]
Description=Chain Mart System

[Service]
EnvironmentFile=/etc/environment
Type=forking
WorkingDirectory=/home/ubuntu
ExecStart=sudo -b /usr/bin/dotnet /home/ubuntu/lib/ChainMart.dll deploy
Restart=always

[Install]
WantedBy=multi-user.target
